# template/profile.d/tmpl_zz-eda-modulepath.sh
# EDA Modulefiles Path Configuration
export MODULEPATH="${MODULE_ROOT}/common:${MODULEPATH:-}"
export MODULEPATH="${MODULE_ROOT}/cadence:${MODULEPATH:-}"
export MODULEPATH="${MODULE_ROOT}/synopsys:${MODULEPATH:-}"
export MODULEPATH="${MODULE_ROOT}/mentor:${MODULEPATH:-}"
export MODULEPATH="${MODULE_ROOT}/other:${MODULEPATH:-}"
